# com.djailla.ssh_server/source/rc.local
#!/bin/sh -e
#
# rc.local
#

case $1 in
  start)
    service ssh start
    ;;
  stop)
    service ssh stop
    ;;
esac
